# docker/crontab
# m h  dom mon dow   command
0 5  * * *  /home/admin/es-outfitter/cron.sh
